# rc.local
#!/bin/sh -e

for governor in /sys/devices/system/cpu/cpufreq/policy*/scaling_governor
do
	echo schedutil > $governor
done

echo 0 > /proc/sys/kernel/sched_util_clamp_min_rt_default

# Throttle all tasks in the system from here onward
/home/qyousef/utils/uclampreset

# Turbo boost is energy enefficient, disable it
#echo 1 > /sys/devices/system/cpu/intel_pstate/no_turbo
